# launchagent/com.local.lark-reaction-as-user.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>REPLACE_LABEL</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_REPO/scripts/lark-reaction-as-user</string>
    <string>--preference</string>
    <string>REPLACE_REPO/preference.yaml</string>
  </array>
  <key>WorkingDirectory</key>
  <string>REPLACE_REPO</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>REPLACE_HOME/.local/state/lark-reaction-as-user/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_HOME/.local/state/lark-reaction-as-user/logs/launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>REPLACE_HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>REPLACE_HOME</string>
  </dict>
</dict>
</plist>
